# provision/plebian-os-firstboot.service
# plebian-os-firstboot.service — provision Plebian-OS on the first boot after a
# regular Debian install: pull deps + the pleb/kilix repos from GitHub and bring
# up the Pleb session. Installed + enabled by the preseed's late_command (or by
# bootstrap.sh). Runs once, then disables itself.
[Unit]
Description=Plebian-OS first-boot provisioning (pull kilix + pleb, set up the Pleb session)
Documentation=https://github.com/itsmygithubacct/plebian-os
# Needs the network (clones from GitHub) and a settled base system.
After=network-online.target systemd-user-sessions.service
Wants=network-online.target
# Hold the login greeter until provisioning finishes, so the very first boot
# never flashes a bare "Default Xsession" (Pleb isn't installed yet) at the
# user — they'd log into that, hit the fallback, and think it's broken. This is
# ORDERING ONLY (not Requires): if provisioning fails or times out, the greeter
# still comes up rather than the boot deadlocking. On every later boot the
# Condition below skips this unit (and it has disabled itself), so the ordering
# is a no-op and normal boots are never delayed.
Before=display-manager.service
# Don't run if provisioning already completed.
ConditionPathExists=!/var/lib/plebian-os/provisioned

[Service]
Type=oneshot
RemainAfterExit=yes
# --kiosk is appended here by the preseed when autologin was requested.
ExecStart=/usr/local/sbin/plebian-os-provision
# Mark done and stop pulling on every boot (success only).
ExecStartPost=/usr/bin/install -Dm644 /dev/null /var/lib/plebian-os/provisioned
ExecStartPost=/bin/systemctl disable plebian-os-firstboot.service
# Provisioning can take a while on a slow link (apt + two clones + engine).
TimeoutStartSec=1800
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
