# provision/plebian-os-firstboot.service
# plebian-os-firstboot.service — provision Plebian-OS on the first boot after a
# regular Debian install: pull deps + the pleb/kilix repos from GitHub and bring
# up the Pleb session. Installed + enabled by the preseed's late_command (or by
# bootstrap.sh). Runs once, then disables itself.
[Unit]
Description=Plebian-OS first-boot provisioning (pull kilix + pleb, set up the Pleb session)
Documentation=https://github.com/itsmygithubacct/plebian-os
# Needs the network (clones from GitHub) and a settled base system.
After=network-online.target systemd-user-sessions.service
Wants=network-online.target
# Hold the login greeter until provisioning finishes, so the very first boot
# never flashes a bare "Default Xsession" (Pleb isn't installed yet) at the
# user — they'd log into that, hit the fallback, and think it's broken. This is
# ORDERING ONLY (not Requires): if provisioning fails or times out, the greeter
# still comes up rather than the boot deadlocking. On every later boot the
# Condition below skips this unit (and it has disabled itself), so the ordering
# is a no-op and normal boots are never delayed.
Before=display-manager.service
# Don't run if provisioning already completed.
ConditionPathExists=!/var/lib/plebian-os/provisioned
# First boot depends on network, mirrors, and GitHub. Three one-hour attempts,
# separated by one minute, all fit inside this four-hour accounting window. A
# fourth start is therefore refused even when every attempt reaches its timeout;
# unlike a window shorter than TimeoutStartSec, retries are genuinely bounded.
StartLimitIntervalSec=4h
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
# Optional build-time provisioning options. The provisioner honours
# PLEBIAN_OS_DESKTOP / PLEBIAN_OS_KIOSK / PLEBIAN_OS_USER; a preseed (e.g. the
# one build_vm_image.py generates) can drop this file to pick the session mode
# without editing the unit. The leading '-' makes it optional.
EnvironmentFile=-/etc/default/plebian-os
# A persistent counter makes the three-attempt bound survive reboots. When it
# is exhausted ExecCondition skips the unit without marking success; the display
# manager can start, and an operator can inspect the journal/fix the network.
ExecCondition=/usr/local/sbin/plebian-os-firstboot-attempt check
ExecStartPre=/usr/local/sbin/plebian-os-firstboot-attempt begin
# Clear any temporary sudoers grant a previously killed attempt may have left
# behind (SIGKILL can't run the provisioner's own cleanup trap).
ExecStartPre=-/bin/rm -f /etc/sudoers.d/plebian-os-provision
ExecStart=/usr/local/sbin/plebian-os-provision
# Defense in depth: cleanup runs after success, ordinary failure, timeout, and
# signal-driven stops even if the shell's own EXIT/signal trap is bypassed.
ExecStopPost=-/bin/rm -f /etc/sudoers.d/plebian-os-provision
# Mark done and stop pulling on every boot (success only).
ExecStartPost=/usr/bin/install -Dm644 /dev/null /var/lib/plebian-os/provisioned
ExecStartPost=/usr/local/sbin/plebian-os-firstboot-attempt success
ExecStartPost=/bin/systemctl disable plebian-os-firstboot.service
Restart=on-failure
RestartSec=60s
# Provisioning can take a while on a slow link, especially now that the kilix
# clickable-chrome fork is built during first boot instead of accepting the
# prebuilt fallback engine as final.
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
